# docs/deploy/supervisor/laravel-worker-exportacao.conf
; Supervisor — Laravel Queue Worker: exportação PDF de Empresas
;
; Copiar para:
;   /etc/supervisor/conf.d/laravel-worker-exportacao.conf
;
; Ajuste /var/www/facigam e user=www-data conforme o servidor.

[program:laravel-worker-exportacao]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/facigam/artisan queue:work --queue=empresas-exportacao --sleep=1 --tries=1 --timeout=900 --max-time=3600
directory=/var/www/facigam
user=www-data
autostart=true
autorestart=true
startsecs=1
startretries=3
numprocs=1
redirect_stderr=false
stdout_logfile=/var/log/supervisor/laravel-worker-exportacao.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/laravel-worker-exportacao.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
stopwaitsecs=930
stopsignal=TERM
killasgroup=true
stopasgroup=true
environment=APP_ENV="production",LARAVEL_QUEUE="empresas-exportacao"
